# dashboard.service
[Unit]
Description=AI Portfolio Manager Dashboard
After=network.target

[Service]
User=root
WorkingDirectory=/opt/portfolio_manager
Environment="PATH=/opt/portfolio_manager/venv/bin"
# Using gunicorn for the Flask app
ExecStart=/opt/portfolio_manager/venv/bin/gunicorn -b 0.0.0.0:8080 dashboard:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
